# scripts/validator/systemd/cathedral-validator.service
[Unit]
Description=Cathedral Validator Service
Documentation=https://github.com/one-covenant/cathedral
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cathedral
Environment=RUST_LOG=debug
ExecStart=/opt/cathedral/cathedral-validator --config /opt/cathedral/config/validator.toml start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cathedral-validator

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/cathedral/data
ReadWritePaths=/opt/cathedral/logs
ProtectHome=yes
RemoveIPC=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
