# scripts/provision/configs/systemd/cathedral-miner.service
[Unit]
Description=Cathedral Miner Service
Documentation=https://github.com/spacejar/cathedral
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/var/lib/cathedral

[Service]
Type=simple
User=cathedral
Group=cathedral
WorkingDirectory=/var/lib/cathedral/miner
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
ExecStart=/usr/local/bin/cathedral-miner --config /etc/cathedral/miner.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cathedral-miner

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/cathedral /var/log/cathedral
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
