# orchestrator/ansible/roles/disk_cleanup/templates/disk-cleanup.service.j2
[Unit]
Description=K3s Disk Cleanup Service
After=k3s.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup-disk.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=disk-cleanup

[Install]
WantedBy=multi-user.target
