# scripts/miner/systemd/cathedral-miner.service
[Unit]
Description=Cathedral Miner Service
Documentation=https://github.com/one-covenant/cathedral
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cathedral
Environment=RUST_LOG=debug
ExecStart=/opt/cathedral/cathedral-miner --config /opt/cathedral/config/miner.toml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cathedral-miner

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/cathedral/data
ReadWritePaths=/opt/cathedral/logs
ReadWritePaths=/root/.ssh
ProtectHome=yes
RemoveIPC=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
